# main.cron
00 03 * * * python /home/demipt/yupi/main.py